# roles/os_hardening/files/anvil-upgrade-notify.service
[Unit]
Description=Anvil — ntfy-Meldung nach unattended-upgrades

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/anvil-upgrade-notify
Nice=10
